# conf/celery/celery_conf.conf
[program:celery]
directory = /app/
command=/usr/local/bin/celery -A bitwatch worker -l info
user=root